# runit/service/kopia/run
#!/bin/sh
exec snooze -H '*' -M '/30' "$(dirname "$0")/backup.sh"
