# runit/service/syncthing/run
#!/bin/sh

exec \
  syncthing \
    serve \
      --no-default-folder
